# systemd/pg_backup.service
[Unit]
Description=Postgresql Backup
Requires=postgresql.service
After=postgresql.service

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/root
ExecStart=/bin/bash /opt/pg_backup/pg_backup.sh
